# jogoweb.service
[Unit]
Description=jogoweb

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/home/ubuntu/jogoweb

[Install]
WantedBy=multi-user.target
